# .config/systemd/user/mpd-log-playing.service
[Unit]
Description=Log MPD played songs

[Service]
ExecStart=/home/pineman/code/bin/mpd-log-playing.sh

[Install]
WantedBy=default.target
