# .root/etc/systemd/system/backup.service
[Service]
ExecStart=/home/pineman/.local/bin/backup.sh

[Install]
WantedBy=multi-user.target
